# build/Tiny-OP-Lienol/files/etc/crontabs/root


*/10 * * * * /etc/BanIp.sh 2>&1 >> /tmp/dropBrute.log
*/30 * * * * sh /root/my_pppoe.sh  restart1
